# supervisor/supervisord.d/chatserver.conf
[program:chatserver]
command=node /usr/local/chatserver/chatserver-wss.js
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/tmp/chatserver.log
